# examples/llm-queue-proxy.service
[Unit]
Description=LLM Queue Proxy Service
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/llm-queue-proxy
EnvironmentFile=/etc/llm-queue-proxy.env
WorkingDirectory=/var/lib/llm-queue-proxy
Restart=on-failure
RestartSec=3
LimitNOFILE=65536
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
