# examples/llm-queue-proxy.service
[Unit]
Description=LLM Queue Proxy Service
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/llm-queue-proxy
WorkingDirectory=/opt/llm-queue-proxy
EnvironmentFile=/etc/llm-queue-proxy.env
Restart=on-failure
RestartSec=3
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
